# middleware/static/config/systemd/system/middleware.service
[Unit]
Description=Unicorn Daemon for middleware
After=network.target

[Service]
User=pi
Group=www-data
WorkingDirectory=/home/pi/way_box_app_v2/middleware
EnvironmentFile=/home/pi/env
ExecStart=/usr/local/bin/gunicorn --workers 2 --bind unix:/home/pi/Way-way_box_app_v2/middleware/middleware.sock wsgi:app

ExecReload=/bin/kill -s HUP $MAINPID
ExecStop=/bin/kill -s TERM $MAINPID


[Install]
WantedBy=multi-user.target
