# middleware/static/config/profile.d/way-box.sh
export API_HOST="api.way-connect.com"
export API_KEY=""
export API_SECRET=""
